# api/root/etc/services.d/api/run
#!/usr/bin/with-contenv sh

cd /data

if [[ ${ENV} == "dev" ]]; then
  reload="--reload"
fi

exec uvicorn opa.main:app --host 0.0.0.0 --log-level ${OPA_LOGLEVEL:-info} ${reload}

