# api/root/etc/services.d/worker/run
#!/usr/bin/with-contenv sh

[[ "${DONTSTART}" == "worker" || "${MODE}" == "api" || "${MODE}" == "manual" ]] && ( /usr/bin/s6-stopdisable worker && exit 0 )

cd /data

exec sleep 1000

