# api/root/etc/services.d/worker/run
#!/usr/bin/with-contenv sh

[[ "${MODE}" == "worker" ]] || ( /usr/bin/s6-stopdisable worker && exit 0 )

cd /data

exec celery ${CELERY_COMMAND:-worker} -A opa.main ${CELERY_PARAMS}
